# csharp/Interfold.Bootstrapper/Phases/SystemdTemplates/interfold.service
[Unit]
Description=Interfold self-hosted backend (docker-compose)
Documentation=https://github.com/Azyyyyyy/octocon
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{OUTPUT_DIR}}
ExecStart=/usr/bin/docker compose -f {{COMPOSE_FILE}} up -d
ExecStop=/usr/bin/docker compose -f {{COMPOSE_FILE}} down
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
